# systemd/crashdummy.service.in
@GENERATED_SCRIPT_FILE_WARNING@

[Unit]
Description=CrashDummy Service
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
ExecStart=@prefix@/sbin/crashdummy
PrivateTmp=true
Restart=no

[Install]
#WantedBy=multi-user.target
